# board/rockchip/rk3036/fs-overlay/etc/profile.d/RkEnv.sh
export HOME=/data
export PATH=$PATH:$HOME:$HOME/bin

# auto run /data/RkLunch.sh
source /data/RkLunch.sh
