# board/rockchip/rk3036/fs-overlay/etc/profile.d/RkEnv.sh
export PATH=$PATH:/data

# auto run /data/RkLunch.sh
source /data/RkLunch.sh
